# shell/demo.crontab
30 21 * * * /usr/local/apache/bin/apachectl restart
